# crontab
*/1 * * * * python manage.py bot >> cron_log.log 2>&1
